# rootfs/etc/services.d/sickchill/run
#!/usr/bin/with-contenv ash

echo "Starting SickChill"
exec s6-setuidgid abc /usr/bin/python /app/SickBeard.py -q --nolaunch --datadir=/config --port 8081