# rootfs/etc/services.d/sickchill/run
#!/usr/bin/with-contenv bash
echo "STARTING SERVICE: sickchill"
/bin/s6-setuidgid abc /app/SickBeard.py --datadir=/config --port 8081 2>&1
